# roles/desktop/templates/syndaemon.service
[Unit]
Description=syndaemon

[Service]
ExecStart=/usr/bin/syndaemon -R -K -i 1
